# cfg/netlabd.conf
[program:netlabd]
command=/usr/bin/netlabd
redirect_stderr=true
stdout_logfile=/var/log/supervisor-netlabd.log
