# cfg/netmgr.conf
[program:netmgr]
command=/usr/bin/netmgr
redirect_stderr=true
stdout_logfile=/var/log/supervisor-netmgr.log
