# cfg/umlabd.conf
[program:umlabd]
command=/usr/sbin/umlabd
redirect_stderr=true
stdout_logfile=/var/log/supervisor-umlabd.log
